# etc/systemd/user/updates-check.timer
[Unit]
Description=Count how many updates are available

[Timer]
OnCalendar=minute
Unit=checkupdates.service

[Install]
WantedBy=timers.target
